# scripts/cryptoedy-subscription-check.timer
[Unit]
Description=Run CryptoEdy subscription check daily at 02:00 UTC

[Timer]
OnCalendar=*-*-* 02:00:00 UTC
# If the system was off at 02:00, run as soon as it comes back up
Persistent=true

[Install]
WantedBy=timers.target
